# docker/etc/supervisord/conf.d/start_emulators.conf
[program:start_4_emulators]
command=/session-appium/docker/start_4_emus.sh
autorestart=false