# docker/etc/supervisor/conf.d/start_emulator_rebuild_state.conf
[program:start_emu_for_state]
command=/session-appium/docker/start_emu_for_state.sh
autorestart=false